# radio-web.service
[Unit]
Description=Radio Web UI
After=network-online.target mpd.service
Wants=network-online.target

[Service]
Type=simple
User=radio
Group=radio
ExecStart=/usr/bin/python3 /usr/local/bin/radio_web.py
Restart=always
RestartSec=2
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/home/radio

[Install]
WantedBy=multi-user.target
