# services/radio-update-stations.service
[Unit]
Description=Update radio stations.yaml from GitHub

[Service]
Type=oneshot
User=radio
ExecStart=/home/radio/radio/scripts/update-stations.sh
